# systemd/news-fetcher.service
[Unit]
Description=News MCP Feed Fetcher
After=network.target

[Service]
Type=simple
User=news-mcp
Group=news-mcp
WorkingDirectory=/opt/news-mcp
Environment=PYTHONPATH=/opt/news-mcp
ExecStart=/opt/news-mcp/venv/bin/python jobs/scheduler.py
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/news-mcp

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=news-fetcher

[Install]
WantedBy=multi-user.target