# systemd/news-api.service
[Unit]
Description=News MCP API Server
After=network.target

[Service]
Type=simple
User=news-mcp
Group=news-mcp
WorkingDirectory=/opt/news-mcp
Environment=PYTHONPATH=/opt/news-mcp
ExecStart=/opt/news-mcp/venv/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/news-mcp

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=news-api

[Install]
WantedBy=multi-user.target